# systemd/scc.service
[Unit]
Description=Scrapyard Command Center UI (Flask)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.config/scc/current
EnvironmentFile=%h/.config/scc/.env
ExecStart=/bin/sh -c '%h/.config/scc/current/.venv/bin/flask --app scc_core.ui run --host="${SCC_UI_BIND:-127.0.0.1}" --port=8081'
Restart=on-failure

[Install]
WantedBy=default.target
